# src/rc.local
#!/bin/dash
# -*- shell-script -*-

# This file is executed in the end of booting
# into multi-user. If this file is not executable,
# it will be sourced into rc.multi instead.

dir=/etc/rc.local.d
if [ -d "$dir" ]; then
    for file in $(ls -1 "$dir"); do
	file = "$dir/$file"
	if [ -f "$file" -a -x "$file" ]; then
	    "$file"
	fi
    done
fi

